# extra/upload.cron
10 0 * * * ( cd && /path/to/python3 extra/upload.py --host ftp://hostname --user username:password ) >> logs/cron/darc-upload.log 2>&1
